# shortcut-fe/shortcut-fe/files/etc/uci-defaults/50-luci-firewall-shortcut-fe
#!/bin/sh

[ -f "/etc/config/ucitrack" ] && {
uci -q batch <<-EOF >/dev/null
	delete ucitrack.@firewall[1]
	add ucitrack firewall
	set ucitrack.@firewall[1].init=shortcut-fe
	commit ucitrack
EOF
}

exit 0
